# old/scripts/base/conf/supervisord-mysqld.conf
[program:mysqld]
command=/usr/local/bin/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true