# scripts/base/conf/supervisord-beanstalkd.conf
[program:beanstalkd]
command=/usr/bin/beanstalkd -l 127.0.0.1 -p 11300 -b /var/lib/beanstalkd
numprocs=1
autostart=true
autorestart=true